# systemd/aoe-bus-watchdog.service
[Unit]
Description=aoe-bus watchdog — auto-nudge stuck aoe Claude sessions
After=graphical-session.target

[Service]
Type=simple
# Use the wrapper installed at ~/.local/bin/aoe-bus so the repo path can change.
ExecStart=%h/.local/bin/aoe-bus watchdog --interval 30 --cooldown 60
Restart=always
RestartSec=10
# Inherit the user's DBus / DISPLAY env so notify-send works from inside.
# (User services on systemd 240+ get these via systemctl --user import-environment)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aoe-bus-watchdog

[Install]
WantedBy=default.target
